# heartbeatd.rc
#!/bin/ksh
#

daemon="/usr/local/sbin/heartbeatd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
